# deploy/schoolbell.service
[Unit]
Description=SchoolBell Gunicorn
After=network.target

[Service]
User=schoolbell
Group=schoolbell
WorkingDirectory=/home/schoolbell/app
EnvironmentFile=/home/schoolbell/app/.env
ExecStart=/home/schoolbell/app/venv/bin/gunicorn wsgi:app \
    --workers 2 \
    --threads 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /home/schoolbell/app/logs/access.log \
    --error-logfile  /home/schoolbell/app/logs/error.log
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
